# cookie/ctyun-login.service
[Unit]
Description=CTYun Login Web Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ctyun
Environment="PATH=/opt/ctyun/venv/bin"
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/ctyun/venv/bin/gunicorn -c /opt/ctyun/gunicorn.conf.py wsgi:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
